# root/etc/s6-overlay/s6-rc.d/init-abs-tracked-mariadb-init/up
/etc/s6-overlay/s6-rc.d/init-abs-tracked-mariadb-init/run
