# root/defaults/crontabs/root
0 3 * * * /usr/local/bin/abs-tracked-backup >>/config/log/backup.log 2>&1
